# templates/etcd.service.erb
[Unit]
Description=Etcd
After=docker.service
Requires=docker.service

[Service]
ExecStartPre=-<%= @docker_path %> rm -f etcd_container
ExecStartPre=<%= @docker_path %> run --rm -v /opt/bin:/opt/bin cellos/systemd-docker
ExecStart=/opt/bin/systemd-docker --cgroups name=systemd run --rm -v /var/etcd-data:/var/etcd-data \
-p 4001:4001 -p 2380:2380 -p 2379:2379 --name etcd_container <%= @registry %>:<%= @version %> \
 -data-dir /var/etcd-data \
 -name etcd<%= @idx %> \
 -advertise-client-urls http://<%= @address %>:2379,http://<%= @address %>:4001 \
 -listen-client-urls http://0.0.0.0:2379,http://0.0.0.0:4001 \
 -initial-advertise-peer-urls http://<%= @address %>:2380 \
 -listen-peer-urls http://0.0.0.0:2380 \
 -initial-cluster-token etcd-cluster-1 \
 -initial-cluster <%= @peers.sort.each_with_index.map{ |peer, index| "etcd#{index}=http://#{peer}:2380" }.join(',') %> \
 -initial-cluster-state new
 
ExecStopPost=-<%= @docker_path %> rm -f etcd_container
Restart=always
RestartSec=10
Type=notify
NotifyAccess=all
TimeoutStartSec=120
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
